# deploy/smart-water-watchdog.service
[Unit]
Description=Taproot - Watchdog Alert
After=network-online.target

[Service]
Type=oneshot
User=jason
WorkingDirectory=/home/jason/taproot
ExecStart=/usr/bin/node src/watchdog.js
EnvironmentFile=/home/jason/.taproot/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=taproot-watchdog
