# deploy/smart-water-summary.service
[Unit]
Description=Smart Water System - Daily Summary Email
After=network-online.target

[Service]
Type=oneshot
User=jason
WorkingDirectory=/home/jason/smart-water
ExecStart=/usr/bin/node src/summary.js
EnvironmentFile=/home/jason/.smart-water/.env
TimeoutStartSec=60
StandardOutput=journal
StandardError=journal
SyslogIdentifier=smart-water-summary
