# deploy/taproot-summary.service
[Unit]
Description=Taproot - Daily Summary Email
After=network-online.target

[Service]
Type=oneshot
User=jason
WorkingDirectory=/home/jason/taproot
ExecStart=/usr/bin/node src/summary.js
EnvironmentFile=/home/jason/.taproot/.env
TimeoutStartSec=60
StandardOutput=journal
StandardError=journal
SyslogIdentifier=taproot-summary
