# deploy/taproot-briefing.service
[Unit]
Description=Taproot - Weekly Intelligence Briefing
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/taproot
ExecStart=/usr/bin/node src/briefing-runner.js
EnvironmentFile=/root/.taproot/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=taproot-briefing
