# deploy/taproot-watchdog.timer
[Unit]
Description=Taproot - Watchdog (runs at 2am daily)

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
