# deploy/taproot-web.service
[Unit]
Description=Taproot - Web UI
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/taproot
ExecStart=/usr/bin/node src/web.js
EnvironmentFile=/root/.taproot/.env
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=taproot-web

[Install]
WantedBy=multi-user.target
